# alpine-build/guirlande-service.initd
#!/sbin/openrc-run

depend()
{
  need soft-pwm-service
}

start()
{
  /usr/sbin/guirlande-service start
}

stop()
{
  /usr/sbin/guirlande-service stop
}
